# supervisord.conf
[supervisord]
nodaemon=true

[program:aspnetdebug]
command=dnx -p project.json aspnet_debug.Server